# contrib/init/procyond.conf
description "Procyon Core Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env PROCYOND_BIN="/usr/bin/procyond"
env PROCYOND_USER="procyon"
env PROCYOND_GROUP="procyon"
env PROCYOND_PIDDIR="/var/run/procyond"
# upstart can't handle variables constructed with other variables
env PROCYOND_PIDFILE="/var/run/procyond/procyond.pid"
env PROCYOND_CONFIGFILE="/etc/procyon/procyon.conf"
env PROCYOND_DATADIR="/var/lib/procyond"

expect fork

respawn
respawn limit 5 120
kill timeout 60

pre-start script
    # this will catch non-existent config files
    # procyond will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$PROCYOND_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run procyond."
        echo "The setting must appear in $PROCYOND_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following "
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"Procyon Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$PROCYOND_PIDDIR"
    chmod 0755 "$PROCYOND_PIDDIR"
    chown $PROCYOND_USER:$PROCYOND_GROUP "$PROCYOND_PIDDIR"
    chown $PROCYOND_USER:$PROCYOND_GROUP "$PROCYOND_CONFIGFILE"
    chmod 0660 "$PROCYOND_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$PROCYOND_PIDFILE" \
    --chuid $PROCYOND_USER:$PROCYOND_GROUP \
    --exec "$PROCYOND_BIN" \
    -- \
    -pid="$PROCYOND_PIDFILE" \
    -conf="$PROCYOND_CONFIGFILE" \
    -datadir="$PROCYOND_DATADIR" \
    -disablewallet \
    -daemon

